# an-website.ini
[program:an-website]
command=/bin/pypy3 /home/user-placeholder/an-website/main.py
directory=/home/user-placeholder/an-website/
autostart=true
process_name=an-website
user=user-placeholder
