# an-website.ini
[program:an-website]
command=/home/user-placeholder/an-website/run.sh
directory=/home/user-placeholder/an-website/
autostart=true
process_name=an-website
user=user-placeholder
stopasgroup=true
killasgroup=true