# an_website.ini
[program:an_website]
command=/home/user-placeholder/an-website/venv/bin/pypy3 -m an_website
directory=/home/user-placeholder/an-website
user=user-placeholder
stopasgroup=true
killasgroup=true
