# an_website.ini
[program:an_website]
command=/home/user-placeholder/an-website/venv/bin/python3 -m an_website
directory=/home/user-placeholder/an-website
user=user-placeholder
stopasgroup=true
stopsignal=INT
environment=HOME="/home/user-placeholder",USER="user-placeholder"
