# an_website.ini
[program:an_website]
command=<home>/an-website/venv/bin/python3 -m an_website
directory=<home>/an-website
user=<user>
stopasgroup=true
stopsignal=INT
environment=LANG="<lang>",LC_ALL="<lang>",HOME="<home>",USER="<user>"
