# an_website.ini
[program:an_website]
command=/home/user-placeholder/an-website/venv/bin/python3 -m an_website
directory=/home/user-placeholder/an-website
user=<user>
stopasgroup=true
stopsignal=INT
environment=LANG="<lang>",LC_ALL="<lang>",HOME="<home>",USER="<user>"
