# an_website.ini
[program:an_website]
command=<home>/an-website/venv/bin/python3 -m an_website
directory=<home>/an-website
autorestart=true
stopsignal=INT
user=<user>
environment=HOME="<home>",LANG="<lang>",LC_ALL="<lang>",SHELL="/bin/sh",USER="<user>"
