# an_website.ini
[program:an_website]
command=/home/user-placeholder/an_website/venv/bin/pypy3 /home/user-placeholder/an_website/main.py
directory=/home/user-placeholder/an_website/
autostart=true
process_name=an_website
user=user-placeholder
stopasgroup=true
killasgroup=true
