# example-configurations/supervisord/an.ini
[program:an_website]
command=/opt/an-website/venv/bin/python3 -m an_website
directory=/opt/an-website
autorestart=true
stopsignal=INT
user=josh
environment=HOME="/home/josh",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8",SHELL="/bin/sh",USER="josh"
