# example-configurations/supervisord/an_website.ini
[program:an_website]
command=/opt/an-website/venv/bin/an-website
directory=/opt/an-website
startsecs=10
autorestart=true
stopwaitsecs=30
killasgroup=true
user=an-website
environment=LANG="de_DE.utf8"
