# example-configurations/supervisord/an_website.ini
[program:an_website]
command=/opt/an-website/venv/bin/python3 -m an_website
directory=/opt/an-website
autorestart=true
startsecs=10
user=<user>
environment=HOME="<home>",LANG="<lang>",LC_ALL="<lang>",SHELL="/bin/sh",USER="<user>"
